# test/services/josuke_2_test.service
[Unit]
Description = weshhhhhh
After       = 

[Service]
ExecStart   = PWD=/home/mkd/work/logchest
Type        = 
Environment = 
